# bin/think-keycodes
#!/bin/bash
# Copyright © 2012 Jim Turner <jturner314@gmail.com>
# Copyright © 2012-2013 Martin Ueding <dev@martin-ueding.de>
# Licensed under The GNU Public License Version 2 (or later)

### BEGIN INIT INFO
# Provides:
# Required-Start: $all
# Required-Stop:
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Sets keycodes for ThinkPad X220
### END INIT INFO

TEXTDOMAIN=think-rotate

case "$1" in
    start|restart|force-reload)
        setkeycodes 6e 109 6d 104 69 28 6b 01 6c 120
        ;;

    stop|status)
        echo "Action $1 is not implemented"
        ;;
esac
