# 81-thinkpad-dock.rules
KERNEL=="dock.0", ACTION=="change", RUN+="/usr/bin/think-dock"

# vim: ft=udevrules
