# 81-thinkpad-dock.rules
KERNEL=="dock.0", ATTR{docked}=="0", RUN+="/usr/bin/think-dock off"
KERNEL=="dock.0", ATTR{docked}=="1", RUN+="/usr/bin/think-dock on"

# vim: ft=udevrules
