# 81-thinkpad-dock.rules
KERNEL=="dock.0", ATTR{docked}=="0", RUN+="/usr/bin/think-rotate off"
KERNEL=="dock.0", ATTR{docked}=="1", RUN+="/usr/bin/think-rotate on"

# vim: ft=udevrules
