# 81-thinkpad-dock.rules
# Copyright © 2012-2013 Martin Ueding <dev@martin-ueding.de>

KERNEL=="dock.0", ATTR{docked}=="0", RUN+="/usr/bin/logger -t think-dock Event caught: off"
KERNEL=="dock.0", ATTR{docked}=="1", RUN+="/usr/bin/logger -t think-dock Event caught: on"
KERNEL=="dock.0", ATTR{docked}=="0", RUN+="/usr/bin/think-dock-hook off"
KERNEL=="dock.0", ATTR{docked}=="1", RUN+="/usr/bin/think-dock-hook on"

KERNEL=="dock.1", ATTR{docked}=="0", RUN+="/usr/bin/logger -t think-dock Event caught: off"
KERNEL=="dock.1", ATTR{docked}=="1", RUN+="/usr/bin/logger -t think-dock Event caught: on"
KERNEL=="dock.1", ATTR{docked}=="0", RUN+="/usr/bin/think-dock-hook off"
KERNEL=="dock.1", ATTR{docked}=="1", RUN+="/usr/bin/think-dock-hook on"

KERNEL=="dock.2", ATTR{docked}=="0", RUN+="/usr/bin/logger -t think-dock Event caught: off"
KERNEL=="dock.2", ATTR{docked}=="1", RUN+="/usr/bin/logger -t think-dock Event caught: on"
KERNEL=="dock.2", ATTR{docked}=="0", RUN+="/usr/bin/think-dock-hook off"
KERNEL=="dock.2", ATTR{docked}=="1", RUN+="/usr/bin/think-dock-hook on"

# vim: ft=udevrules
